# miles-external/package/librespot/librespot.service
[Unit]
Description=Librespot Spotify Connect Client
Documentation=https://github.com/librespot-org/librespot
After=network-online.target chrony.service
Wants=network-online.target
Requires=chrony.service

[Service]
Type=simple
User=librespot
Group=librespot
EnvironmentFile=/etc/default/librespot
ExecStartPre=/bin/mkdir -p "$LIBRESPOT_CACHE" "$LIBRESPOT_SYSTEM_CACHE"
ExecStartPre=/bin/chown librespot:librespot "$LIBRESPOT_CACHE" "$LIBRESPOT_SYSTEM_CACHE"
ExecStart=/usr/bin/librespot \
  --name "$LIBRESPOT_NAME" \
  --device-type "$LIBRESPOT_DEVICE_TYPE" \
  --system-cache "$LIBRESPOT_SYSTEM_CACHE" \
  --cache "$LIBRESPOT_CACHE" \
  --volume-ctrl "$LIBRESPOT_VOLUME_CTRL" \
  --bitrate "$LIBRESPOT_BITRATE" \
  --cache-size-limit "$LIBRESPOT_CACHE_SIZE_LIMIT"
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
